# variant-alpine/etc/apk/protected_paths.d/lbu.list
+home/lncm
+etc/init.d/lncm
+etc/init.d/iotwifi
+etc/init.d/avahi-daemon
