# etc/init.d/avahi-daemon
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

extra_started_commands="reload"

depend() {
	before netmount nfsmount
	use net
	need dbus hostname
	after iotwifi
}

start() {
	ebegin "Starting avahi-daemon"
	/usr/sbin/avahi-daemon -D
	eend $?
}

stop() {
	ebegin "Stopping avahi-daemon"
	/usr/sbin/avahi-daemon -k
	eend $?
}

reload() {
	ebegin "Reloading avahi-daemon"
	/usr/sbin/avahi-daemon -r
	eend $?
}
