# etc/init.d/docker-compose
#!/sbin/openrc-run

# docker-compose init script

depend() {
  need net sshd docker lncm-online
}

start() {
  ebegin "Starting docker-compose"
  
  # Lnd
  if [ ! -d /media/important/important/lnd ]; then
    echo "Error: lnd directory not found"
    exit 1
  fi
    if [ ! -f /media/important/important/lnd/lnd.conf ]; then
    echo "Error: lnd configuration file not found"
    exit 1
  fi
  
  # Bitcoind
  if [ ! -d /media/archive/archive/bitcoin ]; then
    echo "Error: bitcoind directory not found"
    exit 1
  fi
  if [ ! -f /media/archive/archive/bitcoin/bitcoin.conf ]; then
    echo "Error: bitcoin configuration not found"
    exit 1
  fi
  
  # Docker compose
  if [ ! -d /home/lncm/compose ]; then
    echo "Error: docker compose directory not found!"
    exit 1
  fi
  if [ ! -f /home/lncm/compose/docker-compose.yml ]; then
    echo "Error: docker compose configuration not found!"
    exit 1
  fi

  # All good
  cd /home/lncm/compose || exit 1
  docker-compose up -d

  eend $?
}

stop() {
  ebegin "Stopping docker compose"
  
  # Invoke either Noma or Mana
  if [[ -f /usr/bin/noma ]]; then
    /usr/bin/noma stop
  else
    # Invoke Mana if noma doesn't exist
    if [[ -f /usr/local/sbin/mana ]]; then
      /usr/local/sbin/mana stop
    fi
  fi 
  # If neither noma or mana exists just kill docker-compose
  cd /home/lncm/compose || exit 1
  docker-compose down

  eend $?
}
