# etc/init.d/docker-compose
#!/sbin/openrc-run

# docker-compose init script

depend() {
  need net sshd docker lncm-online
}

start() {
  ebegin "Starting docker-compose"

  cd /home/lncm/compose
  docker-compose up -d

  eend $?
}

stop() {
  ebegin "Stopping docker compose"
  # TODO: block until bitcoind has shut down cleanly
  /usr/local/sbin/mana stop
  cd /home/lncm/compose
  docker-compose down

  eend $?
}
