# etc/init.d/iotwifi
#!/sbin/openrc-run

depend() {
	need net docker
}

start() {
	sleep 2 && \
	docker load --input /media/mmcblk0p1/iotwifi.tar && \
	docker run --name=iotwifi -d --restart=unless-stopped --privileged --net host \
		-v /etc/iotwifi:/cfg \
	        -v /etc/wpa_supplicant:/etc/wpa_supplicant \
		cjimti/iotwifi
}
