# etc/init.d/lncm
#!/sbin/openrc-run

# LNCM-Box installation script

depend() {
  need net
  after ntpd sshd
}

start() {
  ebegin "Starting LNCM-Box installation"

  echo "Logging to: /var/log/box-install.log"
  # This script will only run once, for persistent installation to SD
  ./home/lncm/box-install.sh 2>&1 | tee /var/log/box-install.log

  eend $?
}
