# etc/init.d/lncm-online
#!/sbin/openrc-run

# internet connection test

depend() {
  need net sshd
}

start() {
  ebegin "Checking internet connectivity"

  # Fetch current mirror URL
  #MIR="$(cat /etc/apk/repositories | sed -n 2p)"

  # Check if main mirror is reachable
  net_check() {
      ONLINE=$(echo -e 'GET http://dl-cdn.alpinelinux.org/alpine/ HTTP/1.0\n\n' | nc dl-cdn.alpinelinux.org 80 > /dev/null 2>&1; echo $?)
  }

  ONLINE=1
  COUNT=0

  while [ $ONLINE -eq 1 ] && [ $COUNT -lt 10 ]
  do
      echo "Checking for internet connection for 10 seconds..."
      net_check
      COUNT=$((COUNT + 1))
      sleep 1
  done

  if [ $ONLINE -eq 0 ]; then
    echo "Connected to internet"
    exit 0
  else
    echo "Not connected to internet.. please update wpa_supplicant"
    echo "and then reboot to continue the setup process"
    exit 1
  fi

  eend $?
}
