# etc/init.d/lncm-post
#!/sbin/openrc-run
#
# Copyright 2018 LNCM contributors
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#    http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

## LNCM post-installation script
## To be run on persistently installed SD card

# TODO: add system-wide is_online dependency

depend() {
    need net sshd
}

start() {
  ebegin "Starting LNCM post-installation"
  # This script will only run once, to complete post-installation

  if [ -f /media/mmcblk0p1/installed ]; then
    echo "Let apk cache live on persistent volume"
    /sbin/setup-apkcache /var/cache/apk

    echo "Install docker-compose and dependencies for config files"
    # Install dependencies for docker compose and for generating config files
    /sbin/apk update && \
    /sbin/apk add py3-pip curl

    /usr/bin/pip3 install --upgrade pip
    /usr/bin/pip3 install docker-compose
    
    # USB storage setup
    /usr/bin/python3 /usr/local/sbin/lncm-usb.py
    # Verify USB setup succeeded
    if [ ! -d /media/archive ]; then
      echo 'USB setup failed'
      exit 1
    elif [ ! -d /media/important ]; then
      echo 'USB setup failed'
      exit 1
    elif [ ! -d /media/volatile ]; then
      echo 'USB setup failed'
      exit 1
    fi

    # Generate RPC auth
    echo "Generate RPC auth from bitcoin repo"
    curl "https://raw.githubusercontent.com/bitcoin/bitcoin/master/share/rpcauth/rpcauth.py" 2>/dev/null 1>rpcauth.py
    if [ -f ./rpcauth.py ]; then
      chmod 700 rpcauth.py
      ./rpcauth.py lncm | tee /media/important/generated.txt | head -2 | tail -1 > /media/important/rpcauth.txt
      cat generated.txt | tail -1 > /media/important/rpcpass.txt
      rm rpcauth.py /media/important/generated.txt
    else
      echo "rpcauth not generated - could not reach server"
      exit 1
    fi

    if [ $(grep -c rpcauth /media/important/rpcauth.txt ) == 0 ]; then
      echo "Can't generate credentials - exiting!"
      exit 1
    fi
    
    # Apply RPC auth
    if [ -f /media/important/rpcauth.txt ]; then
      # Replace strings with what was generated
      RPCAUTHLINE=`cat /media/important/rpcauth.txt`
      GENERATEDPW=`cat /media/important/rpcpass.txt`
      /bin/sed -i "s/GENERATEDRPCAUTH/$RPCAUTHLINE/g;"  /home/lncm/bitcoin/bitcoin.conf                        
      /bin/sed -i "s/GENERATEDRPCAUTH/$RPCAUTHLINE/g;"  /home/lncm/litecoin/litecoin.conf
      /bin/sed -i "s/RANDOMUSER/lncm/g;"  /home/lncm/lnd/lnd.conf
      /bin/sed -i "s/RANDOMPASS/$GENERATEDPW/g;"  /home/lncm/lnd/lnd.conf
      /bin/sed -i "s/RANDOMUSER/lncm/g;"  /home/lncm/lightning/config
      /bin/sed -i "s/RANDOMPASS/$GENERATEDPW/g;"  /home/lncm/lightning/config
      /bin/sed -i "s/RANDOMUSER/lncm/g;"  /home/lncm/compose/docker-compose.yml
      /bin/sed -i "s/RANDOMPASS/$GENERATEDPW/g;"  /home/lncm/compose/docker-compose.yml
    else
      echo "Error generating config files"
      exit 1                    
    fi
    
    # Install Tor
    echo "Installing and starting TOR (and waiting about 10 seconds)"
    /sbin/apk update && /sbin/apk add tor
    /sbin/service tor start
    
    CHECKTOR_HOSTNAME=1
    # Check to see if Tor generated directory exists
    while [ ! -d /var/lib/tor/lightning ] && [ $CHECKTOR_HOSTNAME -lt 10 ]; 
    do 
      echo "Tor generated directory doesn't exist yet.. waiting";
      let CHECKTOR_HOSTNAME=CHECKTOR_HOSTNAME+1;
      sleep 1;
    done
    
    # Now check to see if hostname exists
    if [ -f /var/lib/tor/lightning/hostname ]; then
      echo "Configuring lnd to have tor as a URI available"
      /bin/sed -i "s/; externalip=/externalip=$(cat /var/lib/tor/lightning/hostname)/g; " /home/lncm/lnd/lnd.conf
    else
      echo "Could not find the tor generated files.. exiting"
      exit 1
    fi

    # Backup apkovl.tar.gz to USB storage
    /sbin/lbu pkg /media/important

    echo "Remove lncm-post from boot"
    /sbin/rc-update del lncm-post default

    exit 0
  else
    echo "LNCM persistent installation not found!"
    exit 1	
  fi

  eend $?
}
