# etc/init.d/lncm-post
#!/sbin/openrc-run

## LNCM post-installation script
## To be run on persistently installed SD card

depend() {
  need net lncm-online
  after sshd ntpd
}

start() {
  ebegin "Starting LNCM post-installation"
  # This script will only run once, to complete post-installation

  check_installed() {
    if [ -f /media/mmcblk0p1/installed ]; then
      echo 0
      return 0
    else
      echo 1
      return 1
    fi
  }

  move_cache() {
    echo "Let apk cache live on persistent volume"
    if [ -d /media/mmcblk0p1/cache ]; then
      cp -v /media/mmcblk0p1/cache/* /var/cache/apk/
    fi
    /sbin/setup-apkcache /var/cache/apk
  }

  install_python3() {
    if ! [ -f /media/mmcblk0p1/nopython ]; then
      echo "Install Python 3 and upgrade PIP"
      /sbin/apk update && \
      /sbin/apk add python3 python3-dev

      /usr/bin/pip3 install --upgrade pip
    fi
  }

  install_utils() {
    echo "Installing utils (avahi/dbus/curl/jq)"
    /sbin/apk update && \
    /sbin/apk add curl jq
  }

  install_docker() {
    if ! [ -f /media/mmcblk0p1/nodocker ]; then
      /sbin/apk add docker
    else
      echo "Skipping docker install"
    fi
  }

  install_avahi() {
    if ! [ -f /media/mmcblk0p1/noavahi ]; then
      echo "Installing AVAHI"
      /sbin/apk add dbus avahi
    else
      echo "Skipping avahi install"
    fi
  }

  config_avahi() {
    if ! [ -f /media/mmcblk0p1/noavahi ]; then
      echo "Configuring Avahi"
      /sbin/rc-update add dbus boot || echo "Unable to add dbus to boot runlevels"
      /sbin/rc-update add avahi-daemon boot || echo "Unable to add avahi to boot runlevels"
      echo "Starting DBUS and AVAHI"
      /etc/init.d/dbus start || echo "Failed to start dbus"
      /etc/init.d/avahi-daemon start || echo "Failed to start avahi"
    fi
  }

  config_docker() {
    if ! [ -f /media/mmcblk0p1/nodocker ]; then
      echo "Configuring docker"
      /sbin/rc-update add docker default
      echo "Starting docker"
      /etc/init.d/docker start || echo "Failed to start docker"
    fi
  }

  install_alpine_sdk() {
    echo "Install alpine devtools (used for building)"
    /sbin/apk add alpine-sdk linux-headers libffi-dev
  }

  install_tor() {
    if ! [ -f /media/mmcblk0p1/notor ]; then
      # Install Tor
      echo "Installing and starting tor (and waiting about 10 seconds)"
      /sbin/apk update && /sbin/apk add tor
      /sbin/service tor start
    else
      echo "Skipping TOR install"
    fi
  }

  enable_tor() {
    if ! [ -f /media/mmcblk0p1/notor ]; then
      echo "Add tor to default runlevels"
      /sbin/rc-update add tor default
    else
      echo "Skipping tor enable"
    fi
  }

  disable_post() {
    echo "Remove lncm-post from boot"
    /sbin/rc-update del lncm-post default
  }

  main() {
    is_installed=$(check_installed)

    if [ ! "$is_installed" -eq "0" ]; then
      echo "Error: LNCM installation not found!"
      exit 1
    fi

    move_cache
    install_utils
    install_docker
    install_avahi
    config_docker
    config_avahi
    install_alpine_sdk
    install_python3
    install_tor
    enable_tor
    disable_post
  }

  main

  eend $?
}
