# etc/init.d/lncm-online
#!/sbin/openrc-run

depend() {
  need net sshd
}

start() {
  ebegin "Checking internet connectivity"

  # Fetch current mirror URL
  #MIR="$(cat /etc/apk/repositories | sed -n 2p)"

  # Check if main mirror is reachable
  net_check() {
      ONLINE=$(echo -e 'GET http://dl-cdn.alpinelinux.org/alpine/ HTTP/1.0\n\n' | nc dl-cdn.alpinelinux.org 80 > /dev/null 2>&1; echo $?)
  }

  ONLINE=1

  while [ $ONLINE -eq 1 ]
  do
      echo "Checking for internet connection..."
      net_check
      sleep 1
  done

  if [ $ONLINE -eq 0 ]; then
    echo "Connected to internet"
    exit 0
  fi

  eend $?
}
