# etc/init.d/nginx
#!/sbin/openrc-run

depend() {
	need net docker iotwifi
}

start() {
	sleep 2 && \
	docker load --input /media/mmcblk0p1/nginx.tar.gz && \
	docker run --name=nginx -d --restart=unless-stopped --net host \
      -v /etc/nginx:/etc/nginx \
      -v /var/log/nginx:/var/log/nginx \
      -v /home/lncm/public_html:/www \
		nginx:1.14.2-alpine
}
