# etc/init.d/noma
#!/sbin/openrc-run

# start up noma services

depend() {
	need net sshd docker
}

start() {
	ebegin "Starting NOMA Services"
	cd /media/noma/compose/neutrino || exit 1
	docker-compose up -d

	noma lnd autounlock || exit 1

	eend $?
}

stop() {
	noma stop || exit 1

	# Attempt to force stop
	cd /media/noma/compose/neutrino || exit 1
	docker-compose down

	eend $?
}
