# etc/apk/protected_paths.d/lbu.list
+etc/motd
+etc/init.d/lncm
+etc/init.d/lncm-post
+etc/init.d/lncm-online
+etc/init.d/docker-compose
+etc/tor/torrc
+etc/init.d/avahi-daemon
+etc/init.d/wpa_supplicant
+etc/init.d/noma
+etc/init.d/sshd
+etc/init.d/wpa_supplicant
+home/lncm/bitcoin
+home/lncm/litecoin
+home/lncm/lnd
+home/lncm/lightning
+home/lncm/compose
+home/lncm/README.md
+home/lncm/crontab
+home/lncm/public_html
+home/lncm/public_html/wifi
+home/lncm/public_html/wifi/index.html
+home/lncm/nginx
+home/lncm/nginx/conf.d
+home/lncm/.ssh/authorized_key
+etc/ssh/sshd_config
+usr/local/bin/setup-disk.sh
+usr/local/sbin
+usr/local/sbin/lncm-getinfo.py
