# etc/apk/protected_paths.d/lbu.list
+etc/tor/torrc
+etc/init.d/lncm
+etc/init.d/lncm-post
+etc/init.d/lncm-online
+etc/init.d/docker-compose
+etc/init.d/iotwifi
+etc/init.d/avahi-daemon
+etc/init.d/portainer
+etc/init.d/sshd
+etc/init.d/nginx
+home/lncm/bitcoin
+home/lncm/litecoin
+home/lncm/lnd
+home/lncm/lightning
+home/lncm/compose
+home/lncm/README.md
+home/lncm/crontab
+home/lncm/public_html
+home/lncm/public_html/wifi
+home/lncm/public_html/wifi/index.html
+home/lncm/nginx
+home/lncm/nginx/conf.d
+home/lncm/.ssh/authorized_keys
+home/lncm/save_password
+home/lncm/seed.txt
+usr/local/sbin
+usr/local/sbin/mana
+usr/local/sbin/lncm-usb
+usr/local/sbin/lncm-usb.py
+usr/local/sbin/lncm-unlock.py
+usr/local/sbin/check-invoicer-status
+usr/local/sbin/lncm-createwallet.py
+usr/local/sbin/lncm-istorprivkeycreated.py
