# etc/apk/protected_paths.d/lbu.list
+etc/init.d/lncm
+etc/init.d/lncm-post
+etc/init.d/iotwifi
+etc/init.d/avahi-daemon
+etc/init.d/portainer
