# etc/apk/protected_paths.d/lbu.list
+etc/tor/torrc
+etc/init.d/lncm
+etc/init.d/lncm-post
+etc/init.d/lncm-online
+etc/init.d/docker-compose
+etc/init.d/iotwifi
+etc/init.d/avahi-daemon
+etc/init.d/portainer
+home/lncm/bitcoin
+home/lncm/litecoin
+home/lncm/lnd
+home/lncm/lightning
+home/lncm/compose
+home/lncm/README.md
+usr/local/sbin/mana
+usr/local/sbin/lncm-usb
+usr/local/sbin/lncm-usb.py