# etc/apk/protected_paths.d/lbu.list
+etc/motd
+etc/init.d/lncm
+etc/init.d/lncm-post
+etc/init.d/lncm-online
+etc/init.d/docker-compose
+etc/init.d/check-authorized-keys
+etc/tor/torrc
+etc/init.d/avahi-daemon
+etc/init.d/wpa_supplicant
+etc/init.d/noma
+etc/init.d/sshd
+etc/init.d/wpa_supplicant
+etc/runlevels/default/check-authorized-keys
+home/lncm/.ssh
+home/lncm/.ssh/authorized_key
+etc/ssh/sshd_config
+usr/local/sbin
+usr/local/sbin/lncm-getinfo.py
+usr/local/sbin/make_modloop.sh
+usr/local/sbin/make_cache.sh
+usr/local/sbin/make_container.sh
+usr/local/sbin/make_modloop.sh
+usr/local/sbin/make_upgrade.sh
